# linux/arch/notify-checkupdates.timer
[Unit]
Description=Update Notifier Timer

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1d
Unit=notify-checkupdates.service

[Install]
WantedBy=default.target
