# linux/systemd/notify-checkupdates.service
[Unit]
Description=Update Notifier

[Service]
Type=oneshot
ExecStart=%h/.local/bin/notify-checkupdates
