# linux/systemd/unclutter.service
[Unit]
Description="unclutter"

[Service]
ExecStart=unclutter --exclude-root --ignore-scrolling
Restart=always

[Install]
WantedBy=default.target
