# linux/systemd/xbanish.service
[Unit]
Description="Hide cursor while typing"

[Service]
ExecStart=xbanish

[Install]
WantedBy=default.target
